# scheduler/windows-task.xml
<?xml version="1.0" encoding="UTF-16"?>
<!-- PFMI ML Prediction Scheduler - Windows Task Scheduler -->
<!-- ======================================================= -->
<!-- Runs ML predictions every 30 minutes on Windows -->
<!--
Installation:
1. Update paths in <Command> and <WorkingDirectory> to match your installation
2. Set environment variables in <EnvironmentVariables>
3. Import this task:
   schtasks /create /xml "C:\path\to\PFMI\scheduler\windows-task.xml" /tn "PFMI ML Predictions"
   
Or use Task Scheduler GUI:
1. Open Task Scheduler (taskschd.msc)
2. Import Task (Action > Import Task)
3. Select this XML file
4. Update paths and credentials
5. Test the task manually before enabling
-->
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2024-07-11T00:00:00</Date>
    <Author>PFMI Admin</Author>
    <Description>Runs ML predictions every 30 minutes and writes results to database. Updates equipment health scores, creates alerts, and stores prediction history.</Description>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <Repetition>
        <Interval>PT30M</Interval>
        <StopAtDurationEnd>false</StopAtDurationEnd>
      </Repetition>
      <StartBoundary>2024-01-01T00:00:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT1H</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>C:\Python\python.exe</Command>
      <Arguments>predict_to_db.py</Arguments>
      <WorkingDirectory>C:\path\to\PFMI\apps\ml-service</WorkingDirectory>
    </Exec>
  </Actions>
</Task>

<!-- 
Environment Variables Setup:
You need to set these as system or user environment variables:
1. Open System Properties > Environment Variables
2. Add these variables:
   - SUPABASE_URL=https://your-project.supabase.co
   - SUPABASE_ANON_KEY=your-anon-key-here
   - ML_SERVICE_URL=http://localhost:8000

Or use PowerShell to set them:
[Environment]::SetEnvironmentVariable("SUPABASE_URL", "https://your-project.supabase.co", "User")
[Environment]::SetEnvironmentVariable("SUPABASE_ANON_KEY", "your-key", "User")
[Environment]::SetEnvironmentVariable("ML_SERVICE_URL", "http://localhost:8000", "User")

Test manually in PowerShell:
cd C:\path\to\PFMI\apps\ml-service
python predict_to_db.py --dry-run
-->
